# quoteTip/quoteTip.cron
50 14 * * 1-5 daemon python3 /data/quoteTip/quoteTip.py
